# rootfs/etc/s6-overlay/s6-rc.d/logger/run
#!/usr/bin/env execlineb

exec s6-log -bp 1 n1 s10000 /var/log/bitcoin
