# rootfs/etc/s6-overlay/s6-rc.d/bitcoin-daemon/run
#!/usr/bin/env sh

: ${BLOCKSDIR}
: ${CHAIN}
: ${CONF}
: ${CONNECT}
: ${DATADIR}
: ${DBCACHE}
: ${DEBUG}
: ${DISABLEWALLET}
: ${DISCOVER}
: ${DNSSEED}
: ${LISTEN}
: ${MAXCONNECTIONS}
: ${PRUNE}
: ${RPCALLOWIP}
: ${RPCPASSWORD}
: ${RPCTHREADS}
: ${RPCUSER}
: ${TXINDEX}
: ${UPNP}
: ${WALLETBROADCAST}
: ${WALLETDIR}
: ${WHITELIST}
: ${ZMQPUBRAWBLOCK}
: ${ZMQPUBRAWTX}

ARGS="$ARGS \
    --blocksdir=${BLOCKSDIR} \
    --chain=${CHAIN} \
    --conf=${CONF} \
    --datadir=${DATADIR} \
    --debug=${DEBUG} \
    --printtoconsole=${PRINTTOCONSOLE} \
    --rpcallowip=${RPCALLOWIP} \
    --walletdir=${WALLETDIR} \
    --zmqpubrawblock=${ZMQPUBRAWBLOCK} \
    --zmqpubrawtx=${ZMQPUBRAWTX} \
"

if [ ! -z "$CONNECT" ]; then
    ARGS="$ARGS --connect=${CONNECT}"
fi

if [ ! -z "$DEBUG" ]; then
    ARGS="$ARGS --debug=${DEBUG}"
fi

if [ ! -z "$DISABLEWALLET" ]; then
    ARGS="$ARGS --disablewallet=${DISABLEWALLET}"
fi

if [ ! -z "$DISCOVER" ]; then
    ARGS="$ARGS --discover=${DISCOVER}"
fi

if [ ! -z "$DNSSEED" ]; then
    ARGS="$ARGS --dnsseed=${DNSSEED}"
fi

if [ ! -z "$LISTEN" ]; then
    ARGS="$ARGS --listen=${LISTEN}"
fi

if [ ! -z "$MAXCONNECTIONS" ]; then
    ARGS="$ARGS --maxconnections=${MAXCONNECTIONS}"
fi

if [ ! -z "$PRUNE" ]; then
    ARGS="$ARGS --prune=${PRUNE}"
fi

if [ ! -z "$RPCPASSWORD" ]; then
    ARGS="$ARGS --rpcpassword=${RPCPASSWORD}"
fi

if [ ! -z "$RPCTHREADS" ]; then
    ARGS="$ARGS --rpcthreads=${RPCTHREADS}"
fi

if [ ! -z "$RPCUSER" ]; then
    ARGS="$ARGS --rpcuser=${RPCUSER}"
fi

if [ ! -z "$TXINDEX" ]; then
    ARGS="$ARGS --txindex=${TXINDEX}"
fi

if [ ! -z "$UPNP" ]; then
    ARGS="$ARGS --upnp=${UPNP}"
fi

if [ ! -z "$WALLETBROADCAST" ]; then
    ARGS="$ARGS --walletbroadcast=${WALLETBROADCAST}"
fi

if [ ! -z "$WHITELIST" ]; then
    ARGS="$ARGS --whitelist=${WHITELIST}"
fi

if echo "$(cat /var/log/bitcoin/current && cat /var/log/bitcoin/*.s 2> /dev/null)" | grep -q \
    -e "-reindex" \
; then
    ARGS="$ARGS --reindex"
fi

exec /usr/local/bin/bitcoind $ARGS
